# LifeArcade/scripts/com.google.lifearcade.kiosk.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Unique service identifier -->
    <key>Label</key>
    <string>com.google.lifearcade.kiosk</string>

    <!-- Command to execute -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Applications/LifeArcade.app/Contents/MacOS/launch.sh</string>
    </array>

    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>/Users/specialguestx/Documents/GitHub/SGx_GoogleEmployment/LifeArcade</string>

    <!-- Start on login -->
    <key>RunAtLoad</key>
    <true/>

    <!-- ALWAYS restart if Chrome closes (crash or manual close) -->
    <key>KeepAlive</key>
    <true/>

    <!-- Wait 10 seconds between restarts -->
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <!-- Environment variables (required for Docker) -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
